# udev-rules/90-c14hcopy-mount.rules
# CAM
SUBSYSTEMS=="usb", ENV{ID_FS_LABEL}=="JVCCAM_SD", ACTION=="add", RUN+="/bin/sh -c "TERM=linux echo JVCCAM found. starting copy. > /dev/tty0"", RUN+="/bin/mount '%E{DEVNAME}' /mnt/cam", RUN+="/usr/bin/systemctl start c14hcopy@cam"
# MIC
SUBSYSTEMS=="usb", ENV{ID_FS_LABEL}=="ZOOM", ACTION=="add", RUN+="/bin/sh -c "TERM=linux echo ZOOM MIC found. starting copy. > /dev/tty0"", RUN+="/bin/mount '%E{DEVNAME}' /mnt/mic", RUN+="/usr/bin/systemctl start c14hcopy@mic"
